# backend/xiaou-telemetry-rollup.timer
[Unit]
Description=Run XiaoU Todo telemetry rollup every hour

[Timer]
OnBootSec=3min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
